# supervisor/supervisord.conf
[unix_http_server]
file=/tmp/naviguide-sup.sock
chmod=0700

[supervisord]
logfile=/mnt/efs/spaces/ef014a98-8a1c-4b16-8e06-5d2c5b364d08/527f928c-fa41-40dd-8d61-c174d3e76a01/feature_deploy-ui-fixes/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/tmp/naviguide-sup.pid
nodaemon=false
silent=false
user=cocoapp

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/naviguide-sup.sock

[program:naviguide-api]
command=/usr/bin/python3 main.py
directory=/mnt/efs/spaces/ef014a98-8a1c-4b16-8e06-5d2c5b364d08/527f928c-fa41-40dd-8d61-c174d3e76a01/feature_deploy-ui-fixes/naviguide-api
; API key loaded from .env.naviguide — never commit secrets to git
environment=ANTHROPIC_API_KEY="%(ENV_ANTHROPIC_API_KEY)s"
autostart=true
autorestart=true
startretries=5
startsecs=5
stopasgroup=true
killasgroup=true
stdout_logfile=/mnt/efs/spaces/ef014a98-8a1c-4b16-8e06-5d2c5b364d08/527f928c-fa41-40dd-8d61-c174d3e76a01/feature_deploy-ui-fixes/logs/naviguide-api.log
stderr_logfile=/mnt/efs/spaces/ef014a98-8a1c-4b16-8e06-5d2c5b364d08/527f928c-fa41-40dd-8d61-c174d3e76a01/feature_deploy-ui-fixes/logs/naviguide-api.log
stdout_logfile_maxbytes=10MB
priority=100

[program:naviguide-polar]
command=/usr/bin/python3 -m uvicorn polar_api.main:app --host 0.0.0.0 --port 8004
directory=/mnt/efs/spaces/ef014a98-8a1c-4b16-8e06-5d2c5b364d08/527f928c-fa41-40dd-8d61-c174d3e76a01/feature_deploy-ui-fixes/naviguide_workspace
environment=PYTHONPATH="/mnt/efs/spaces/ef014a98-8a1c-4b16-8e06-5d2c5b364d08/527f928c-fa41-40dd-8d61-c174d3e76a01/feature_deploy-ui-fixes:/mnt/efs/spaces/ef014a98-8a1c-4b16-8e06-5d2c5b364d08/527f928c-fa41-40dd-8d61-c174d3e76a01/feature_deploy-ui-fixes/polar_agent"
autostart=true
autorestart=true
startretries=5
startsecs=5
stopasgroup=true
killasgroup=true
stdout_logfile=/mnt/efs/spaces/ef014a98-8a1c-4b16-8e06-5d2c5b364d08/527f928c-fa41-40dd-8d61-c174d3e76a01/feature_deploy-ui-fixes/logs/polar-api.log
stderr_logfile=/mnt/efs/spaces/ef014a98-8a1c-4b16-8e06-5d2c5b364d08/527f928c-fa41-40dd-8d61-c174d3e76a01/feature_deploy-ui-fixes/logs/polar-api.log
stdout_logfile_maxbytes=10MB
priority=200

[program:naviguide-proxy]
command=/usr/bin/python3 proxy_server.py
directory=/mnt/efs/spaces/ef014a98-8a1c-4b16-8e06-5d2c5b364d08/527f928c-fa41-40dd-8d61-c174d3e76a01/feature_deploy-ui-fixes
environment=PORT="3016",NAVIGUIDE_LOG_DIR="/mnt/efs/spaces/ef014a98-8a1c-4b16-8e06-5d2c5b364d08/527f928c-fa41-40dd-8d61-c174d3e76a01/feature_deploy-ui-fixes/logs"
autostart=true
autorestart=true
startretries=5
startsecs=8
stopasgroup=true
killasgroup=true
stdout_logfile=/mnt/efs/spaces/ef014a98-8a1c-4b16-8e06-5d2c5b364d08/527f928c-fa41-40dd-8d61-c174d3e76a01/feature_deploy-ui-fixes/logs/proxy_server.log
stderr_logfile=/mnt/efs/spaces/ef014a98-8a1c-4b16-8e06-5d2c5b364d08/527f928c-fa41-40dd-8d61-c174d3e76a01/feature_deploy-ui-fixes/logs/proxy_server.log
stdout_logfile_maxbytes=10MB
priority=300
